# etc/supervisord.d/transcode.ini
[program:transcode]
command = node /srv/transcode.js
environment = USER=nobody,HOME=/tmp
autostart = true
autorestart = true
stdout_events_enabled=true
stderr_events_enabled=true
